# systemd/ip_checker.service
[Unit]
Description=Run ip-checker

[Service]
EnvironmentFile=/path/to/secure/env_file
ExecStart=/path/to/your/script.py
